# brados-sec.service
[Unit]
Description=BradOS Security Daemon — policy enforcement, file integrity, threat scans
Documentation=https://github.com/Architect-Brad/BradOS
After=network.target

[Service]
Type=simple
ExecStart=PYTHON_BIN BRADOS_HOME/brados.py --daemon
WorkingDirectory=BRADOS_HOME
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=default.target
